# etc/init/nginx.conf
description "nginx_front"
author "Me"
start on filesystem and started docker
stop on runlevel [!2345]
respawn
script
  #/usr/bin/docker run -v /data/containers/nginx_front/config/nginx: 
end script
